# platform/clounix/sonic-platform-modules-pegatron/fn8031-bnf/service/fn8031_bnf-platform-init.service
[Unit]
Description=Pegastron fn8031-bnf Platform initialization service
Before=network.target
DefaultDependencies=no

[Service]
Type=simple
ExecStart=/usr/local/bin/pegaProcess/main.py install
ExecStop=/usr/local/bin/pegatron_fn8032_bnf_util.py uninstall
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
